# resources/cav_nd_agent/cavisson/monitors/bin/cmon
# chkconfig: 2345 98 02
### BEGIN INIT INFO
# Provides:             cmon
# Required-Start:       
# Required-Stop:        
# Should-Start:         
# Should-Stop:          
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    cmon
### END INIT INFO

#-----------------------------------------------------------------------------------
# Name          : cmon
# Purpose       : To start cmon on system reboot with user provided in cmon.env 
#                        cmon <-----------------
#                         |                     |
#                         V                     |If cmon stop due to any reason, start cmon
#                     cmon_control              ^
#                         |                     |
#                         V                     |
#                     nsu_check_cmon ----->-----
#
# Date          : Sat Oct 17 2015
# Author(s)     : Manish Kr. Mishra
# Modification\ : 
# Date
#------------------------------------------------------------------------------------

OS_NAME=`uname`

#Set CAV_MON_HOME using ~ or /opt/cavisson/monitors
#Commenting it so that it can set according to the environment to support multiple cmon on same machine

#Setting CAV_MON_HOME as nsu_server_admin send .. as a CAV_MON_HOME as it dosen't have CAV_MON_HOME in case of running without topology 
if [ "X$CAV_MON_HOME" == "X.." ];then
  PWD=`pwd`
  CAV_MON_HOME=`dirname $PWD`
fi

# If CAV_MON_HOME is Not set or bin dir does not exist, then set it
if [ "X$CAV_MON_HOME" = "X" -o ! -d $CAV_MON_HOME/bin ];then
  # If CavMon installation path is users home directory
  # Below we have replaced Tilde(~) with $HOME because ~ in not working in SunOS.
  if [ -d $HOME/cavisson/monitors ];then
    if [ -d $HOME/cavisson/monitors/bin ];then
     CAV_MON_HOME=$HOME/cavisson/monitors
    else
      echo "CavMon package is not installed in ~/cavisson/monitors directory. Exiting ..."
      exit 255
    fi
  # If some other CavMon installation path
  elif [ ! -d /opt/cavisson/monitors ];then
    if [ "X$OS_NAME" = "XSunOS" ]; then
      cav_mon=`dirname $0|nawk -F"/bin" '{print $1}'`
    else
      cav_mon=`dirname $0|awk -F"/bin" '{print $1}'`
    fi
    if [ -d $cav_mon ]; then
      if [ -d $cav_mon/bin ]; then
        CAV_MON_HOME=$cav_mon
      else
        echo "CavMon package is not installed in $cav_mon directory. Exiting ..."
        exit 255
      fi
    else
      echo "CavMon package is not installed. Exiting ..."
      exit 255
   fi
  # NS/NO
  elif [ -d /opt/cavisson/monitors ];then
     NS_WDIR="/home/netstorm/work"
     CAV_MON_HOME="/opt/cavisson/monitors"
  
     if [ -f $NS_WDIR/etc/netstorm.env ];then
       . $NS_WDIR/etc/netstorm.env
     #May be installed in some customer machine, which does not have NetStorm build so, env file will not present
     #else
     #  echo "$NS_WDIR/etc/netstorm.env is not found. Exiting ..."
     #  exit 255
     fi
  else
    echo "CavMon package is not installed. Exiting ..."
    exit 255
  fi
fi

export CAV_MON_HOME
export OS_NAME

#nf env set
nfenv="$CAV_MON_HOME/bin/nfenv"
if [ -f "$nfenv" ];then
  source $CAV_MON_HOME/bin/nfenv
fi
#End nf env set

#Passing cmon as second argument in cmon_control since in Khols from cmon.service cmon_control is called and need to run cmon in foreground
if [ `id -u` -eq 0 -a "XX$CMON_USER" != "XXroot" ];then
  su $CMON_USER "-c $CAV_MON_HOME/bin/cmon_control $1 cmon"
else
  $CAV_MON_HOME/bin/cmon_control $1 cmon
fi

exit 0
